# conf/supervisord.conf
[supervisord]
logfile=/tmp/supervisord.log
logfile_maxbytes=5MB
pidfile=/tmp/supervisord.pid
user=root
nodaemon=true
stdout_logfile=/dev/stdout

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:etcd] 
command=/usr/bin/etcd --name=master --data-dir=${PATH_ETCD_DATA}/${FILE_ETCD_BD_NAME} --listen-client-urls=http://0.0.0.0:2379 --advertise-client-urls=http://0.0.0.0:2379
stderr_logfile=/dev/stderr
stdout_logfile=/dev/stdout
stdout_events_enabled=true
stderr_events_enabled=true
redirect_stderr=true
priority=0

[program:etcd-script] 
command=/usr/bin/etcdctl mkdir /funceme.br/network && /usr/bin/etcdctl mk /funceme.br/network/config "{ \"Network\": \"${SERVICE_CLUSTER_IP_RANGE}\", \"SubnetLen\": 24, \"Backend\": { \"Type\": \"vxlan\" } }"
autostart=true
autorestart=false
stderr_logfile=/dev/stderr
stdout_logfile=/dev/stdout
stdout_events_enabled=true
stderr_events_enabled=true
redirect_stderr=true
priority=100
